# src/app/service/myservice.conf.in
description     "xyzmyservice"
author          "asem<asem_alla@msn.com>"

start on started network-services
stop on runlevel [016]

console log

exec @SERVICE_EXECUTABLE_PATH@

respawn
respawn limit 5 60
